# templates/systemd/wg-easy.service.j2
# {{ ansible_managed }}
[Unit]
Description=The easiest way to run WireGuard VPN with Web-based Admin UI
After=network.target

[Service]
Environment=DEBUG=Server,Wireguard
Environment=PASSWORD={{ wireguard_easy_password }}
Environment=LISTEN_IP={{ wireguard_easy_listen_ip }}
Environment=PORT={{ wireguard_easy_port }}
Environment=WG_DEFAULT_ADDRESS={{ wireguard_easy_wg_default_address }}
Environment=WG_DEFAULT_DNS={{ wireguard_easy_wg_default_dns }}
Environment=WG_HOST={{ wireguard_easy_wg_host }}
Environment=WG_PORT={{ wireguard_easy_wg_port }}
Environment=WG_PATH={{ wireguard_easy_config_path }}
Environment='WG_POST_UP=iptables -t nat -A POSTROUTING -s {{ wireguard_easy_wg_default_address }} -o {{ wireguard_easy_host_interface }} -j MASQUERADE -m comment --comment "wg-easy rule";\
iptables -A INPUT -p udp -m udp --dport {{ wireguard_easy_wg_port }} -j ACCEPT -m comment --comment "wg-easy rule";\
iptables -A FORWARD -i wg0 -j ACCEPT -m comment --comment "wg-easy rule";\
iptables -A FORWARD -o wg0 -j ACCEPT -m comment --comment "wg-easy rule";\
'

User={{ wireguard_easy_user }}
WorkingDirectory={{ wireguard_easy_bin_path }}/app
ExecStart=node server.js

[Install]
WantedBy=default.target
